# systemd/audio-fix.service
[Unit]
Description=Audio Fix

[Service]
ExecStart=/home/codex/bin/audio-fix

[Install]
WantedBy=multi-user.target
